# ops/systemd/bet-agent-pipeline.service
[Unit]
Description=Bet Agent daily pipeline
Requires=docker.service
After=docker.service bet-agent-compose.service

[Service]
Type=oneshot
WorkingDirectory=/opt/bet-agent
ExecStart=/usr/bin/docker compose run --rm --profile ops bet-agent-pipeline
